# systemd/tempod.service
[Unit]
Description=tempo background scheduler
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -lc '%h/.local/bin/tempod'
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
